# templates/etcd-backup.timer.j2
[Unit]
Description=Run daily backup of etcd

[Timer]
{{ etcd_backup_periodicity }}
Persistent=true

[Install]
WantedBy=timers.target
